# scripts/backup-db.service
[Service]
Type=oneshot
ExecStart=backup-db