# scripts/deploy-btcmap.org.service
[Service]
Type=simple
WorkingDirectory=/etc/btcmap.org
ExecStart=/etc/btcmap.org/deploy.sh