# scripts/fetch-elements.service
[Service]
Type=simple
ExecStart=fetch-elements