# scripts/fetch-elements.service
[Service]
Type=simple
Environment=DISCORD_WEBHOOK_URL={{ discord_webhook_url }}
ExecStart=fetch-elements