# scripts/generate-report.service
[Service]
Type=oneshot
ExecStart=btcmap-api generate-report