# scripts/generate-legacy-nodes.service
[Service]
Type=oneshot
ExecStart=generate-legacy-nodes