# scripts/generate-legacy-nodes.timer
[Timer]
OnCalendar=*-*-* *:5,15,25,35,45,55:00
Persistent=true

[Install]
WantedBy=timers.target