# scripts/report-legacy-nodes.service
[Service]
Type=oneshot
Environment=DISCORD_WEBHOOK_URL={{ discord_webhook_url }}
ExecStart=report-legacy-nodes