# scripts/sync-users.service
[Service]
Type=oneshot
ExecStart=btcmap-api sync-users