# etc/systemd/system/sync-users.timer
[Timer]
OnCalendar=*-*-* 01:05:00
Persistent=true

[Install]
WantedBy=timers.target